# infcore/deploy/systemd/infcore-gateway.service
# infcore gateway - systemd unit для offline РФ-контура.
# ПРЕДВАРИТЕЛЬНО (один раз, от root):
#   useradd -r -s /usr/sbin/nologin infcore          # сервисный пользователь
#   install -d -o infcore -g infcore /opt/infcore/bin /opt/infcore/config
#   # положить бинари в /opt/infcore/bin, конфиг в /opt/infcore/config/gateway.yaml
#   # каталог журнала /var/log/infcore создаёт сам systemd (LogsDirectory=, см. ниже)
# Установка: cp в /etc/systemd/system/, затем systemctl daemon-reload && enable --now.
[Unit]
Description=infcore gateway (OpenAI-совместимый шлюз к локальным LLM)
After=network.target
Documentation=file:/opt/infcore/docs/DEPLOY.md

[Service]
Type=simple
User=infcore
Group=infcore

# Секреты (API-ключи через env:VAR в конфиге) подхватываются отсюда; файл опционален.
EnvironmentFile=-/etc/infcore/gateway.env

ExecStart=/opt/infcore/bin/infcore_gateway /opt/infcore/config/gateway.yaml

# Gateway сам поднимает дочерние llama-server. control-group гарантирует, что при
# остановке/перезапуске юнита будут убиты ВСЕ процессы (включая модели), без сирот.
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

Restart=on-failure
RestartSec=2

# Базовое усиление. /dev оставлен доступным для GPU (NVIDIA); audit-журнал - на запись.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
# systemd создаёт /var/log/infcore с владельцем User= (mode 0750) до старта сервиса -
# отдельный mkdir/chown не нужен; каталог автоматически в ReadWritePaths.
LogsDirectory=infcore
LogsDirectoryMode=0750
LimitNOFILE=65536

# Offline-инвариант на уровне ядра: запрещаем весь egress, кроме loopback.
# Клиенты ходят на gateway - добавьте их подсеть в IPAddressAllow (напр. 10.0.0.0/8).
IPAddressDeny=any
IPAddressAllow=localhost

[Install]
WantedBy=multi-user.target
